# root/etc/s6-overlay/s6-rc.d/init-opencode-server/up
/etc/s6-overlay/s6-rc.d/init-opencode-server/run